# files/etc/rc.local
# Put your custom commands here that should be executed once
# the system init finished. By default this file does nothing.

# 新增：重新挂载根文件系统为可读写状态
mount -o remount,rw /

exit 0
